# packaging/linux/onevox.service
[Unit]
Description=Onevox Speech-to-Text Daemon
Documentation=https://github.com/kssgarcia/onevox
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/onevox daemon --foreground
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Resource limits (optional, can be adjusted per system)
MemoryMax=2G
CPUQuota=80%

[Install]
WantedBy=default.target
